# cards.service
[Unit]
Description=uWSGI instance to server Card Game Page
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/var/www/andrewbounds/CardGamePage
Environment="PATH=/var/www/andrewbounds/AndrewBoundsEnv/bin"
ExecStart=/var/www/andrewbounds/AndrewBoundsEnv/bin/uwsgi --ini CardGame.ini

[Install]
WantedBy=multi-user.target